# package/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        systemctl daemon-reload >/dev/null 2>&1 || true

        if [ -z "$2" ]; then
            systemctl enable fseventbridge.service >/dev/null 2>&1 || true
            echo "fseventbridge installed. Enable with: systemctl enable --now fseventbridge"
        else
            echo "Upgrading fseventbridge..."
            systemctl try-restart fseventbridge.service >/dev/null 2>&1 || true
        fi
        ;;
esac

exit 0
