# package/prerm
#!/bin/sh
set -e

case "$1" in
    remove|deconfigure)
        # 卸载或降级前停止并禁用服务
        systemctl stop FsEventBridge.service >/dev/null 2>&1 || true
        systemctl disable FsEventBridge.service >/dev/null 2>&1 || true
        ;;
esac

exit 0
